# conf.d/flask.conf
[program:flask]
command=python3 /app/app.py
autorestart=true
